# templates/systemd/prometheus-varnish-exporter.service.j2
[Unit]
Description=Prometheus Varnish Exporter
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/prometheus-varnish-exporter

[Install]
WantedBy=multi-user.target
