# templates/systemd/prometheus-haproxy-exporter.service.j2
[Unit]
Description=Prometheus Haproxy Exporter
After=network.target
 
[Service]
Type=simple
ExecStart=/usr/local/bin/prometheus-haproxy-exporter
 
[Install]
WantedBy=multi-user.target
