# templates/systemd/prometheus-memcached-exporter.service.j2
[Unit]
Description=Prometheus Memcached Exporter
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/prometheus-memcached-exporter

[Install]
WantedBy=multi-user.target
