# templates/systemd/prometheus_phpfpm_exporter.service.j2
[Unit]
Description=Prometheus PHP-FPM Exporter
After=network.target
 
[Service]
Type=simple
ExecStart=/usr/local/bin/prometheus-phpfpm-exporter
 
[Install]
WantedBy=multi-user.target
